# rootfs/etc/systemd/system/hass-install.service
[Unit]
Description=Install HASS Supervised

[Service]
ExecStart=/entrypoint.sh

[Install]
WantedBy=multi-user.target