# rootfs/etc/systemd/system/hass-install.service
[Unit]
Description=Install HASS Supervised
After=docker.service

[Service]
ExecStart=/entrypoint.sh

[Install]
WantedBy=multi-user.target